# packaging/general/sie-dns-sensor/etc/cron.d/sie-trim-spool
*/5 * * * * root /usr/lib/sie/sie-trim-spool
